# deploy/com.regimetrader.rebalance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  regime-trader cross-sectional book DAILY rebalance (launchd LaunchAgent).
  Runs `main.py --rebalance --execute` every weekday (Mon-Fri) ~22:30 CEST, after the
  US close. Daily cadence with split timescales so the bot is attentive every day
  without churning a slow signal: the momentum SELECTION re-ranks only on the first run
  of a new calendar month; intra-month daily runs keep the month's names and only
  re-scale total gross exposure to today's market vol via the vol-target overlay
  (de-risk on vol spikes, re-risk on calm). The open-order guard (drop_open_order_symbols)
  makes a same-day re-run a no-op. PAPER account only; real money remains BLOCKED until
  the gate passes.

  Install:
    cp deploy/com.regimetrader.rebalance.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.regimetrader.rebalance.plist
  Stop:
    launchctl unload ~/Library/LaunchAgents/com.regimetrader.rebalance.plist
  Run now (test):
    launchctl start com.regimetrader.rebalance

  NOT loaded by this session: enabling it (and later switching to live submission)
  is Pablo's gated step, after a supervised paper verification.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.regimetrader.rebalance</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/pablomiguelgonzalezprado/AIOS/projects/regime-trader/.venv/bin/python</string>
    <string>main.py</string>
    <string>--rebalance</string>
    <string>--execute</string>
  </array>

  <!-- T0.4 determinism (R-4): multithreaded BLAS reductions are order-
       nondeterministic; identical HMM fits diverged at ~5e-13 and EM amplified it.
       Single-threaded numeric libs make fits bit-reproducible. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>OMP_NUM_THREADS</key><string>1</string>
    <key>OPENBLAS_NUM_THREADS</key><string>1</string>
    <key>MKL_NUM_THREADS</key><string>1</string>
    <key>VECLIB_MAXIMUM_THREADS</key><string>1</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/pablomiguelgonzalezprado/AIOS/projects/regime-trader</string>

  <!-- Every weekday (Weekday 1-5 = Mon-Fri) at 22:30 local time, after the US close. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/pablomiguelgonzalezprado/AIOS/projects/regime-trader/logs/rebalance.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/pablomiguelgonzalezprado/AIOS/projects/regime-trader/logs/rebalance.err.log</string>
</dict>
</plist>
