# deploy/com.regimetrader.runonce.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  regime-trader daily paper cycle (launchd LaunchAgent).
  Runs `main.py --run-once` once per US trading day, ~22:15 CEST
  (after the 16:00 ET close), acting on the freshly-closed daily bar; bracket
  orders queue and fill at the next open.

  Install:
    cp deploy/com.regimetrader.runonce.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.regimetrader.runonce.plist
  Stop:
    launchctl unload ~/Library/LaunchAgents/com.regimetrader.runonce.plist
  Run now (test):
    launchctl start com.regimetrader.runonce

  Note: launchd cannot skip US market holidays; on a holiday the cycle re-runs
  on the unchanged daily bar and the C4 delta gate yields no new order (hold).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.regimetrader.runonce</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/pablomiguelgonzalezprado/AIOS/projects/regime-trader/.venv/bin/python</string>
    <string>main.py</string>
    <string>--run-once</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/pablomiguelgonzalezprado/AIOS/projects/regime-trader</string>

  <!-- Weekdays (Mon=1 .. Fri=5) at 22:15 local time. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>15</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>15</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>15</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>15</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>15</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/pablomiguelgonzalezprado/AIOS/projects/regime-trader/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/pablomiguelgonzalezprado/AIOS/projects/regime-trader/logs/launchd.err.log</string>
</dict>
</plist>
